# resources/blattfruit.cron
*/2 * * * * ~/.blattfruit/blattfruit.sh
