# resources/blattfruit.cron
*/5 * * * * ~/.blattfruit/blattfruit.sh
